# config/services/vault-api.service
[Unit]
Description=Today Vault Read API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{deploy_path}}
ExecStart=/usr/bin/node {{deploy_path}}/src/vault-api.js
Restart=on-failure
RestartSec=10
Environment="NODE_ENV=production"
Environment="HOME={{deploy_path}}"
Environment="VAULT_API_PORT=3334"
# Use the same API key as inbox API
EnvironmentFile=-{{deploy_path}}/.inbox-api-key

# Logging
StandardOutput=append:{{deploy_path}}/.data/vault-api.log
StandardError=append:{{deploy_path}}/.data/vault-api.log

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
