# config/services/vault-watcher.service
[Unit]
Description=Today Vault File Watcher
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{deploy_path}}
ExecStart=/usr/bin/node {{deploy_path}}/bin/vault-watcher
Restart=on-failure
RestartSec=10
Environment="NODE_ENV=production"
Environment="HOME={{deploy_path}}"

# Logging
StandardOutput=append:{{deploy_path}}/.data/vault-watcher.log
StandardError=append:{{deploy_path}}/.data/vault-watcher.log

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
