# config/inbox-api.service
[Unit]
Description=Today Inbox Upload API
Documentation=https://github.com/OlderGay-Men/today
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/today
ExecStart=/usr/bin/node /opt/today/src/inbox-api.js
Restart=on-failure
RestartSec=10
Environment="NODE_ENV=production"
Environment="HOME=/opt/today"
Environment="INBOX_API_PORT=3333"
# Generate a random API key on first install, or use existing
EnvironmentFile=-/opt/today/.inbox-api-key

# Logging
StandardOutput=append:/opt/today/.data/inbox-api.log
StandardError=append:/opt/today/.data/inbox-api.log

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target