# config/resilio-sync.service
[Unit]
Description=Resilio Sync
Documentation=https://help.resilio.com
After=network.target

[Service]
Type=simple
User=rslsync
Group=rslsync
RuntimeDirectory=resilio-sync
RuntimeDirectoryMode=0755
ExecStartPre=/bin/mkdir -p /var/run/resilio-sync
ExecStartPre=/bin/chown rslsync:rslsync /var/run/resilio-sync
ExecStart=/usr/bin/rslsync --config /etc/resilio-sync/config.json --nodaemon
Restart=always
RestartSec=5
# Restart daily as backup if scheduler isn't running
RuntimeMaxSec=86400
# Memory limit - restart if exceeds 512MB
MemoryMax=512M

[Install]
WantedBy=multi-user.target
